# deploy/bare-metal/systemd/milnet-ratchet.service
# ==============================================================================
# MILNET SSO — Ratchet Session Manager (Port 9105)
# ==============================================================================
# Implements Double Ratchet session key rotation. Manages session state,
# generates forward-secret session keys, and coordinates with the verifier
# for heartbeat checks and with risk for anomaly detection.
#
# Completely unreachable from the public bastion tier.
# ==============================================================================

[Unit]
Description=MILNET SSO - Ratchet (session key rotation)
Documentation=file:///opt/milnet/docs/ratchet.md
After=network-online.target milnet-pre.service
Wants=network-online.target
Requires=milnet-pre.service
PartOf=milnet.target

[Service]
Type=exec
User=milnet-ratchet
Group=milnet

ExecStart=/opt/milnet/bin/ratchet
Restart=on-failure
RestartSec=5
WatchdogSec=60

# ── Filesystem hardening ──
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/milnet/ratchet
PrivateTmp=yes
PrivateDevices=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectHostname=yes
ProtectClock=yes

# ── Process hardening ──
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
MemoryDenyWriteExecute=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
ProtectProc=invisible
ProcSubset=pid
PrivateUsers=yes

# ── Resource limits ──
# Ratchet holds ephemeral session keys in memory. 256M is generous for
# up to ~100k concurrent sessions.
MemoryMax=256M
CPUQuota=100%
TasksMax=32
LimitNOFILE=65536
LimitMEMLOCK=infinity

# ── Security context ──
SecureBits=noroot
UMask=0077

# ── Environment ──
EnvironmentFile=/etc/milnet/ratchet.env

# ── Logging ──
StandardOutput=journal
StandardError=journal
SyslogIdentifier=milnet-ratchet

[Install]
WantedBy=milnet.target
