# deploy/bare-metal/systemd/milnet-postgres.service
# ==============================================================================
# MILNET SSO — PostgreSQL Database Service
# ==============================================================================
# Drop-in wrapper for system PostgreSQL, ensuring it runs with hardened
# settings appropriate for MILNET. This is a Type=oneshot that validates
# PostgreSQL is properly configured before declaring ready.
#
# PostgreSQL is the highest-value data target in the SSO system. Only
# orchestrator, audit, and admin are permitted database access (enforced
# separately via pg_hba.conf and nftables).
# ==============================================================================

[Unit]
Description=MILNET SSO - PostgreSQL Database (wrapper)
Documentation=man:postgres(1)
After=network-online.target
Wants=network-online.target

# If the system PostgreSQL unit has a different name, adjust here.
# On Debian/Ubuntu: postgresql@15-main.service
# On RHEL/Fedora: postgresql-15.service
Requires=postgresql.service
After=postgresql.service

[Service]
Type=oneshot
RemainAfterExit=yes

# Verify PostgreSQL is listening and the milnet_sso database exists.
# This gate prevents dependent services from starting against a
# misconfigured or uninitialized database.
ExecStart=/bin/bash -c '\
    for i in $(seq 1 30); do \
        if pg_isready -h 127.0.0.1 -p 5432 -U milnet -d milnet_sso -q 2>/dev/null; then \
            echo "MILNET: PostgreSQL is ready"; \
            exit 0; \
        fi; \
        sleep 1; \
    done; \
    echo "MILNET: PostgreSQL failed readiness check after 30s"; \
    exit 1'

# ── Security context ──
# This unit only runs a health check; it needs minimal privileges.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=milnet.target
