# deploy/bare-metal/systemd/milnet-pre.service
# ==============================================================================
# MILNET SSO — Pre-flight Checks
# ==============================================================================
# Runs before any MILNET service starts. Validates:
#   1. Secure Boot chain integrity (if available)
#   2. vTPM PCR measurements against expected baseline
#   3. Binary attestation — SHA-512 hashes of all MILNET binaries
#   4. FIPS mode Known-Answer Tests (KAT) for crypto primitives
#
# If any check fails, this unit fails, and all MILNET services refuse to
# start (they have After/Requires on this unit). This is intentional:
# a tampered boot chain or modified binary means the system is compromised.
# ==============================================================================

[Unit]
Description=MILNET SSO - Pre-flight Security Checks
Documentation=file:///opt/milnet/docs/security-checks.md
DefaultDependencies=no
Before=milnet.target
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Run all security validation scripts in sequence.
# Any failure aborts the entire MILNET stack.
ExecStart=/opt/milnet/security/secure-boot-verify.sh
ExecStart=/opt/milnet/security/vtpm-attest.sh
ExecStart=/opt/milnet/security/binary-attest.sh

# ── Security context ──
# Pre-flight needs read access to /sys/firmware (Secure Boot), /dev/tpmrm0
# (TPM resource manager), and /opt/milnet/bin (binary hashes).
User=root
Group=root
NoNewPrivileges=yes
ProtectHome=yes
PrivateTmp=yes
ProtectKernelModules=yes
ProtectClock=yes
ProtectHostname=yes

[Install]
WantedBy=milnet.target
