# deploy/bare-metal/systemd/milnet-opaque.service
# ==============================================================================
# MILNET SSO — OPAQUE Password Authentication Service (Port 9102)
# ==============================================================================
# Implements the OPAQUE asymmetric PAKE protocol. Handles password
# registration and authentication without ever seeing the plaintext
# password (zero-knowledge proof). Only the orchestrator may communicate
# with this service.
#
# Extremely sensitive: holds OPAQUE server keys. Memory must be locked
# to prevent swap exposure.
# ==============================================================================

[Unit]
Description=MILNET SSO - OPAQUE (password authentication)
Documentation=file:///opt/milnet/docs/opaque.md
After=network-online.target milnet-postgres.service milnet-pre.service
Wants=network-online.target
Requires=milnet-pre.service
BindsTo=milnet-postgres.service
PartOf=milnet.target

[Service]
Type=exec
User=milnet-opaque
Group=milnet

ExecStart=/opt/milnet/bin/opaque
Restart=on-failure
RestartSec=5
WatchdogSec=60

# ── Filesystem hardening ──
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/milnet/opaque
PrivateTmp=yes
PrivateDevices=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectHostname=yes
ProtectClock=yes

# ── Process hardening ──
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
MemoryDenyWriteExecute=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
ProtectProc=invisible
ProcSubset=pid
PrivateUsers=yes

# ── Resource limits ──
# OPAQUE is compute-light but memory-sensitive (key material).
MemoryMax=256M
CPUQuota=100%
TasksMax=32
LimitNOFILE=65536
LimitMEMLOCK=infinity

# ── Security context ──
SecureBits=noroot
UMask=0077

# ── Environment ──
EnvironmentFile=/etc/milnet/opaque.env

# ── Logging ──
StandardOutput=journal
StandardError=journal
SyslogIdentifier=milnet-opaque

[Install]
WantedBy=milnet.target
