# deploy/bare-metal/systemd/milnet-gateway.service
# ==============================================================================
# MILNET SSO — Gateway Service (Port 9100)
# ==============================================================================
# Public-facing bastion service. Accepts puzzle challenges from the internet,
# performs X-Wing KEM key exchange, and forwards authenticated requests to
# the orchestrator. This is the ONLY service reachable from external networks.
#
# Security rationale for every directive is documented inline.
# ==============================================================================

[Unit]
Description=MILNET SSO - Gateway (public bastion)
Documentation=file:///opt/milnet/docs/gateway.md

# Wait for network, database wrapper, and pre-flight checks.
After=network-online.target milnet-postgres.service milnet-pre.service
Wants=network-online.target
Requires=milnet-pre.service
PartOf=milnet.target

[Service]
Type=exec
User=milnet-gateway
Group=milnet

ExecStart=/opt/milnet/bin/gateway
Restart=on-failure
RestartSec=5

# Watchdog: gateway must notify systemd within 60s or it is killed.
# Protects against deadlocks in the puzzle challenge path.
WatchdogSec=60

# ── Filesystem hardening ──
# strict: entire filesystem is read-only except explicit ReadWritePaths.
# Prevents an attacker who achieves code exec from modifying system files.
ProtectSystem=strict
# yes: /home, /root, /run/user are inaccessible. Gateway has no user data.
ProtectHome=yes
# Gateway state directory for rate-limit counters and connection tracking.
ReadWritePaths=/var/lib/milnet/gateway
# Private /tmp prevents cross-service tmp symlink attacks.
PrivateTmp=yes
# No access to physical devices (/dev/sda, /dev/mem, etc.).
PrivateDevices=yes
# Prevent modification of cgroup hierarchy.
ProtectControlGroups=yes
# Prevent loading/unloading kernel modules.
ProtectKernelModules=yes
# Prevent writing to /proc/sys, /sys (sysctl).
ProtectKernelTunables=yes
# Prevent reading kernel ring buffer (dmesg).
ProtectKernelLogs=yes
# Prevent changing system hostname.
ProtectHostname=yes
# Prevent changing system clock (NTP manipulation attacks).
ProtectClock=yes

# ── Process hardening ──
# Irrevocably prevent gaining new privileges via execve (suid, capabilities).
NoNewPrivileges=yes
# Drop ALL capabilities. Gateway binds to port 9100 (>1024), so no
# CAP_NET_BIND_SERVICE needed. If binding to port 443, uncomment below.
CapabilityBoundingSet=
AmbientCapabilities=
# Uncomment if gateway binds to port < 1024 (e.g., 443):
# CapabilityBoundingSet=CAP_NET_BIND_SERVICE
# AmbientCapabilities=CAP_NET_BIND_SERVICE
# Prevent creating executable memory regions (JIT). All MILNET binaries
# are ahead-of-time compiled Rust — no JIT needed.
MemoryDenyWriteExecute=yes
# Lock the execution domain so personality(2) cannot be used to weaken
# security (e.g., switching to 32-bit mode to exploit different syscall ABI).
LockPersonality=yes
# No realtime scheduling — prevents priority inversion DoS.
RestrictRealtime=yes
# Prevent creating setuid/setgid files.
RestrictSUIDSGID=yes
# Prevent creating user/network/pid namespaces (container escape vector).
RestrictNamespaces=yes
# Only allow IPv4, IPv6, and Unix sockets. No netlink, packet, etc.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
# Whitelist only system-service syscalls. Blocks mount, reboot, kexec, etc.
SystemCallFilter=@system-service
# Return EPERM instead of killing the process on denied syscalls.
# Allows the Rust panic handler to produce a useful backtrace.
SystemCallErrorNumber=EPERM
# Completely hide other processes in /proc.
ProtectProc=invisible
ProcSubset=pid
# Map service UID into a private user namespace (defense-in-depth).
PrivateUsers=yes

# ── Resource limits ──
# Gateway handles puzzle challenges which are CPU-intensive by design.
# 512M is sufficient for connection tracking + puzzle state.
MemoryMax=512M
# Allow full CPU for puzzle computation bursts.
CPUQuota=100%
# 64 threads max (connection handler pool).
TasksMax=64
# 65536 file descriptors for high-concurrency connection handling.
LimitNOFILE=65536
# Allow mlock() for pinning cryptographic key material in RAM
# (prevents keys from being swapped to disk).
LimitMEMLOCK=infinity

# ── Security context ──
# Prevent the process from ever gaining root capabilities.
SecureBits=noroot
# Files created by gateway are owner-only (0700 dirs, 0600 files).
UMask=0077

# ── Environment ──
EnvironmentFile=/etc/milnet/gateway.env

# ── Logging ──
# Send stdout/stderr to the journal with a stable identifier.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=milnet-gateway

[Install]
WantedBy=milnet.target
