# deploy/bare-metal/systemd/milnet.target
# ==============================================================================
# MILNET SSO — Service Group Target
# ==============================================================================
# Groups all MILNET SSO services under a single target for unified
# management. Usage:
#
#   systemctl start milnet.target    # Start all MILNET services
#   systemctl stop milnet.target     # Stop all MILNET services
#   systemctl status milnet.target   # Check overall status
#
# Individual services declare PartOf=milnet.target so stopping this
# target propagates to all members.
# ==============================================================================

[Unit]
Description=MILNET SSO - All Services
Documentation=file:///opt/milnet/docs/deployment.md
After=network-online.target
Wants=milnet-pre.service
Wants=milnet-postgres.service
Wants=milnet-audit.service
Wants=milnet-opaque.service
Wants=milnet-tss@1.service milnet-tss@2.service milnet-tss@3.service milnet-tss@4.service milnet-tss@5.service
Wants=milnet-ratchet.service
Wants=milnet-verifier.service
Wants=milnet-orchestrator.service
Wants=milnet-gateway.service
Wants=milnet-admin.service

[Install]
WantedBy=multi-user.target
